# snapraid/services/snapraid-runner.timer
[Unit]
Description="SnapRAID Runner Script Daily 3 AM."

[Timer]
OnCalendar=3:01
Unit=snapraid-runner.service

[Install]
WantedBy=multi-user.target
